# debian/asterisk.postrm
#! /bin/sh

set -e

if [ "$1" = purge ]; then
        userdel -r asterisk 2>/dev/null || true
        rm -fR /var/log/asterisk

fi

# The following bit is placed here manually instead of through the
# <hash>DEBHELPER<hash> replacement.
if [ "$1" = "purge" ] ; then
    update-rc.d asterisk remove >/dev/null
fi

#DEBHELPER#
